# image/runit/php-fpm/run
#!/bin/sh
exec 2>&1
mkdir -p /run/php/
chmod 777 /run/php
exec /usr/sbin/php-fpm7.1 --nodaemonize --fpm-config /etc/php/7.1/fpm/php-fpm.conf
