# bamboo-remote-agent-php55/runit/bamboo-agent/run
#!/bin/sh
exec 2>&1
exec /sbin/setuser bamboo /srv/bamboo-agent-home/bin/bamboo-agent.sh console
